# Buster/chown-sites.service
[Service]
Type=simple
ExecStart=/opt/scripts/chown-sites.sh

[Install]
WantedBy=multi-user.target